# crontab
0 */2 * * *  python /request.py > /proc/1/fd/1 2>/proc/1/fd/2
